# deploy/algo-paper.service
[Unit]
Description=Algo Paper Trader (BurstMomentum strategy on live Binance WS)
After=network.target docker.service algo-api.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/algo-trading
# Default: taker fill model (conservative). Add --maker for optimistic.
ExecStart=/home/ubuntu/algo-trading/.venv/bin/python -m workers.paper_trader
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
